# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Nginx - 反向代理（等待Next.js健康后启动，避免冷启动502）
[program:nginx]
command=/bin/sh -c 'echo "[startup] waiting for Next.js health on 127.0.0.1:3000"; retries=0; until curl -fsS http://127.0.0.1:3000/api/health >/dev/null 2>&1; do retries=$((retries+1)); if [ "$retries" -ge 60 ]; then echo "[startup] Next.js health check timeout after 60s" >&2; exit 1; fi; sleep 1; done; exec nginx -g "daemon off;"'
autostart=true
autorestart=true
startretries=3
startsecs=10
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=12
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=300

# Next.js应用 - Web服务器（优先启动）
[program:nextjs]
command=node server.js
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=12
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=100
user=nextjs
environment=NODE_ENV="production",HOME="/app",OPENCLAW_STATE_DIR="/app/.openclaw",OPENCLAW_WORKSPACE_DIR="/app/.openclaw/workspace",PORT="3000",HOSTNAME="0.0.0.0",REDIS_URL="%(ENV_REDIS_URL)s",QUEUE_SPLIT_BACKGROUND="true",QUEUE_BACKGROUND_WORKER="false",NODE_OPTIONS="--max-old-space-size=%(ENV_NODE_MAX_OLD_SPACE_SIZE_WEB)s"

# 调度器 - 定时任务
[program:scheduler]
command=node dist/scheduler.js
directory=/app
autostart=true
autorestart=true
startretries=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=12
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=300
user=nextjs
environment=NODE_ENV="production",HOME="/app",OPENCLAW_STATE_DIR="/app/.openclaw",OPENCLAW_WORKSPACE_DIR="/app/.openclaw/workspace",TZ="Asia/Shanghai",REDIS_URL="%(ENV_REDIS_URL)s",QUEUE_SPLIT_BACKGROUND="true",QUEUE_BACKGROUND_WORKER="false",NODE_OPTIONS="--max-old-space-size=%(ENV_NODE_MAX_OLD_SPACE_SIZE_SCHEDULER)s"

# Background Queue Worker - 非核心任务（click-farm / url-swap）
# 🔧 修复: 移除 Next.js 健康检查依赖，background-worker 可以独立启动
# 原因: background-worker 只依赖 Redis，不需要等待 Next.js
# 好处: 避免因健康检查超时导致进入 FATAL 状态
[program:background-worker]
command=node dist/background-worker.js
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=10
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=12
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=250
user=nextjs
environment=NODE_ENV="production",HOME="/app",OPENCLAW_STATE_DIR="/app/.openclaw",OPENCLAW_WORKSPACE_DIR="/app/.openclaw/workspace",TZ="Asia/Shanghai",REDIS_URL="%(ENV_REDIS_URL)s",QUEUE_BACKGROUND_WORKER="true",QUEUE_SPLIT_BACKGROUND="true",QUEUE_SKIP_STARTUP_PENDING_REPAIR="true",NODE_OPTIONS="--max-old-space-size=%(ENV_NODE_MAX_OLD_SPACE_SIZE_BACKGROUND_WORKER)s"

# OpenClaw Gateway - 独立 Node22 运行
[program:openclaw-gateway]
command=/usr/local/bin/node22 /app/openclaw/openclaw.mjs gateway
directory=/app/openclaw
autostart=true
autorestart=true
startretries=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=12
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=240
user=nextjs
environment=NODE_ENV="production",HOME="/app",DATABASE_URL="",DATABASE_PATH=":memory:",OPENCLAW_CONFIG_PATH="/app/.openclaw/openclaw.json",OPENCLAW_STATE_DIR="/app/.openclaw",OPENCLAW_WORKSPACE_DIR="/app/.openclaw/workspace",OPENCLAW_BUNDLED_SKILLS_DIR="/app/openclaw/skills",OPENCLAW_LOG_FILE="/proc/self/fd/1",OPENCLAW_CONSOLE_LEVEL="info",OPENCLAW_CONSOLE_STYLE="compact",NODE_OPTIONS="--max-old-space-size=%(ENV_NODE_MAX_OLD_SPACE_SIZE_OPENCLAW)s"

# Python Google Ads Service - 服务账号模式API
[program:python-ads-service]
command=python3 main.py
directory=/app/python-service
autostart=true
autorestart=true
startretries=3
stopsignal=TERM
stopasgroup=true
killasgroup=true
stopwaitsecs=12
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=200
user=nextjs
